# dictation.service
[Unit]
Description=Dictation Service
After=graphical-session.target pulseaudio.service
PartOf=graphical-session.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
KillSignal=SIGTERM
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
WorkingDirectory=@@WORKING_DIR@@
ExecStart=/bin/bash @@EXEC_START@@
Restart=no
RestartSec=10
TimeoutStartSec=30

[Install]
WantedBy=graphical-session.target 
